# webhook-server.service
[Unit]
Description=GitHub Webhook Server
After=network.target

[Service]
User=kali
WorkingDirectory=/home/kali/catty-reminders-app
ExecStart=/usr/bin/python3 /home/kali/catty-reminders-app/webhook_server.py
Restart=always

[Install]
WantedBy=multi-user.target
